# package/sysconfdir/files/profile.d/alias.sh
alias mc='mc -a'
alias sc='setconsole'
alias sr='setconsole -r'

alias xn='killall neutrino'
alias bn='/usr/bin/neutrino'
alias tn='/tmp/neutrino'
alias sn='/etc/init.d/start_neutrino'
