# skel-root/general/files-etc/profile.d/alias.sh
alias mc='mc -a'
alias sc='setconsole'
alias sr='setconsole -r'
